# mac/packages/ToC/守护进程/com.youqu.todesk.session.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>com.youqu.todesk.mac</string>
	</array>
	<key>Disabled</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.youqu.todesk.desktop</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>LoginWindow</string>
		<string>Aqua</string>
	</array>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/ToDesk.app/Contents/Helpers/ToDesk_Session</string>
		<string>-RunAsAgent</string>
		<string>YES</string>
		<string>-Module</string>
		<string>Full</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>WorkingDirectory</key>
	<string>/Applications/ToDesk.app/Contents/MacOS</string>
</dict>
</plist>
